# packages/vm-modules-mounter/vm-modules-mounter.service
[Unit]
Description=Mount kernel modules
DefaultDependencies=no
After=systemd-remount-fs.service
Before=sysinit.target systemd-modules-load.service systemd-udevd.service kmod-static-nodes.service

[Install]
WantedBy=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/vm-modules-mounter.sh
